# update-dns/etc/systemd/system/update-dns.timer
[Unit]
Description=Update DNS A record once every 15 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target